# etc/systemd/exabgp.service
[Unit]
Description=ExaBGP
After=network.target

[Service]
Environment=exabgp_daemon_daemonize=false
ConditionPathExists=/etc/exabgp/exabgp.conf
ExecStart=/usr/sbin/exabgp --folder /etc/exabgp /etc/exabgp/exabgp.conf
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
