# etc/systemd/exabgp.service
[Unit]
Description=ExaBGP
After=network.target

[Service]
Environment=exabgp.daemon.daemonize=false
Environment=exabgp.log.destination=stdout
ExecStart=/usr/sbin/exabgp /etc/exabgp/exabgp.conf

[Install]
WantedBy=multi-user.target
